# server/scanner-backup.service
[Unit]
Description=Scanner DB Backup
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/scanner
ExecStart=/opt/scanner/scripts/backup.sh
User=scanner

[Install]
WantedBy=multi-user.target
