# server/scanner-server.service
[Unit]
Description=물류창고 스캐너 API 서버
After=network.target

[Service]
Type=simple
User=scanner
WorkingDirectory=/opt/scanner/server
EnvironmentFile=/opt/scanner/server/.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
